# deploy/templates/systemd/jiecool-backend.service
# JieCool 后端服务 Systemd 配置
# 版本: v1.0.0
# 描述: 管理 JieCool 后端服务的启动、停止和重启

[Unit]
Description=JieCool Backend Service
Documentation=https://github.com/your-username/JieCool
After=network.target postgresql.service
Wants=postgresql.service
ConditionPathExists=/home/jiecool/server/main

[Service]
Type=simple
User=jiecool
Group=jiecool
WorkingDirectory=/home/jiecool
ExecStart=/home/jiecool/server/main
ExecReload=/bin/kill -HUP $MAINPID

# 重启配置
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# 环境变量
Environment=GO_ENV=production
Environment=GIN_MODE=release

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jiecool-backend

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/jiecool/logs /home/jiecool/uploads /home/jiecool/backups

# 网络配置
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true

# 文件系统权限
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target